# .dotfiles/.config/systemd/user/rclone-onedrive-sync.service
[Unit]
Description=Sync OneDrive to local folder via rclone
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone bisync /home/gaizhe/OneDrive onedrive: \
    --log-file=/home/gaizhe/.cache/rclone/onedrive-sync.log \
    --log-level=INFO
Environment=RCLONE_CONFIG=/home/gaizhe/.config/rclone/rclone.conf

[Install]
WantedBy=default.target
